# roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter

[Service]
Restart=always
ExecStart={{ node_exporter_path }}/node_exporter

[Install]
WantedBy=multi-user.target
