# roles/ceph/templates/ceph_exporter.service.j2
[Unit]
Description=Ceph metrics exporter
After=docker.service

[Service]
Restart=always
ExecStart=/bin/docker run --name=ceph_exporter \
                          -v {{ ceph_config_dir }}:/etc/ceph \
                          -p {{ ceph_exporter_port }}:9148 \
                          {{ ceph_exporter_docker_image }}
ExecStop=/bin/docker stop -t 2 ceph_exporter
ExecStopPost=/bin/docker rm -f ceph_exporter

[Install]
WantedBy=multi-user.target
